# certbot/churaverse-cert-renew.timer
[Unit]

Description=Monthly Churaverse cert renewal timer

[Timer]
OnCalendar= *-*-13 15:00:00
RandomizedDelaySec=1h
Persistent=true
